# Ansible_Roles/Supervisor-Role/templates/long_script.conf.j2
[program:{{ item.conf_filename }}]
command= {{ item.command }}
autostart=true
autorestart=true
stderr_logfile= {{ item.stderr_logfile }}
stdout_logfile= {{ item.stdout_logfile }}
